# server/systemd/palworld-events.timer
[Unit]
Description=Refresh the public Palworld event history every 30 seconds

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=1s
Persistent=true
Unit=palworld-events.service

[Install]
WantedBy=timers.target
