# server/systemd/palworld-stats.timer
[Unit]
Description=Collect Palworld local statistics every minute

[Timer]
OnBootSec=60s
OnUnitActiveSec=60s
AccuracySec=1s
Unit=palworld-stats.service

[Install]
WantedBy=timers.target
