# server/systemd/palworld-events.service
[Unit]
Description=Collect privacy-safe Palworld public events
After=local-fs.target

[Service]
Type=oneshot
User=gaylemon
Group=gaylemon
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
CPUQuota=35%
CPUWeight=25
IOWeight=25
MemoryMax=768M
TimeoutStartSec=300
UMask=0077
ExecStart=/usr/bin/python3 /home/gaylemon/Gaylemon/server/bin/palworld-events-collect.py --backfill-from 2026-07-09T00:00:00-04:00 --backfill-budget 2 --min-backfill-fps 45 --max-backfill-frame-ms 28 --max-backfill-load 5.5 --recent-limit 2000
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
PrivateNetwork=true
ProtectSystem=strict
ProtectHome=read-only
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX
SystemCallArchitectures=native
SystemCallFilter=@system-service
ReadOnlyPaths=/home/gaylemon/Gaylemon/runtime/public-save-snapshot.json /home/gaylemon/Gaylemon/runtime/public-save-bases.json /home/gaylemon/Gaylemon/runtime/save-snapshot-history /home/gaylemon/Gaylemon/runtime/save-bases-history /srv/storage/steam/servers/palworld/stats/stats.json
ReadWritePaths=/home/gaylemon/Gaylemon/runtime
