# server/systemd/palworld-save-snapshot.timer
[Unit]
Description=Check for a fresh Palworld save snapshot every 30 seconds

[Timer]
OnBootSec=45s
OnUnitInactiveSec=30s
AccuracySec=1s
Persistent=true
Unit=palworld-save-snapshot.service

[Install]
WantedBy=timers.target
