# server/systemd/palworld-save-snapshot.timer
[Unit]
Description=Check for a fresh Palworld save snapshot every minute

[Timer]
OnBootSec=60s
OnUnitInactiveSec=60s
AccuracySec=1s
Persistent=true
Unit=palworld-save-snapshot.service

[Install]
WantedBy=timers.target
